# roles/storagebox/files/rc.local
modprobe cifs
echo 0 > /proc/fs/cifs/OplockEnabled
